# files/maintenance/pgsql_vacuum.cron
30 0 * * *  postgres    for db in `/usr/bin/psql --no-align -t -c 'SELECT datname FROM pg_database;' | /bin/grep -Ev '(^template[0-9]$|^postgres$)'`; do /usr/bin/vacuumdb -d $db -f -z | /bin/grep -v VACUUM; done
